# nunatak/resource/nginx-super.ini
[program:nginx]
command = {{nginx}}/sbin/nginx -c {{etc_dir}}/nginx.conf
priority = 10
redirect_stderr = true


